# utils/vct/experiments/hello-world-debian/etc/init.d/confine-experiment
#! /bin/sh

### BEGIN INIT INFO
# Provides:             confine-experiment
# Required-Start:       $syslog
# Required-Stop:        $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         
# Short-Description:    confine-experiment
### END INIT INFO

set -e

#set -x

export PATH="${PATH:+$PATH:}/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin"

case "$1" in
  start)
        echo -n "Starting $DESC: "
        echo "$NAME"
	echo "$(date) started" >> /root/confine-exp.log
	/root/confine-exp.sh /root/confine-exp.log &
        ;;

  stop)
        echo -n "Stopping $DESC: "
        echo "$NAME"
	echo "$(date) stopped (explicitly)" >> /root/confine-exp.log
	killall confine-exp.sh
        ;;

  reload|force-reload|restart)
        $0 stop
        sleep 1
        $0 start
        ;;

  status)
        ps aux | grep confine-exp.sh
        ;;

  *)
        N=/etc/init.d/$NAME
        echo "Usage: $N {start|stop|restart}" >&2
        exit 1
        ;;
esac

exit 0
