# packages/confine/confine-community-lab/files/etc/uci-defaults/disable-firewall.sh
#!/bin/sh
echo "Disabling firewall..."
chmod 000 /etc/init.d/firewall 

