# packages/confine/confine-system/files/etc/init.d/confine-init
#!/bin/sh /etc/rc.common

START=21

start() {
    confine_node_init
}
