# packages/confine/confine-system/files/etc/init.d/confine
#!/bin/sh /etc/rc.common


START=81


start() {

    confine.lib confine_system_check

    confine_sliver_stop all

    if [ "$( uci -q get confine-node.node.state )" = "started" ] ; then
	confine.lib uci_set confine-node.node.state=applied
    fi

    confine_node_setup
}


stop() {

    confine_sliver_stop all

    if [ "$( uci -q get confine-node.node.state )" = "started" ] ; then
	confine.lib uci_set confine-node.node.state=applied
    fi
    
}
