# packages/confine/confine-system/files/etc/init.d/confine
#!/bin/sh /etc/rc.common


START=81


start() {

    confine.lib confine_system_check

    confine_sliver_stop all

    if [ "$( uci -q get confine.node.state )" = "started" ] ; then
	confine.lib uci_set confine.node.state=applied
    fi

    confine_node_enable

    if ! [ "$( uci -q get confine.node.start_daemon )" = "0" ] ; then
	/usr/lib/lua/confine/confine.lua > /dev/null 2>&1 &
    fi

}


stop() {

    if [ -f /var/run/confine/pid ]; then
	kill $(cat /var/run/confine/pid)
    fi

    confine_sliver_stop all

    if [ "$( uci -q get confine.node.state )" = "started" ] ; then
	confine.lib uci_set confine.node.state=applied
    fi
    
}
