# packages/confine/confine-system/files/etc/init.d/confine
#!/bin/sh /etc/rc.common


START=81
STOP=20
FIRSTB="/etc/firstboot"
FIRSTB_BACKUP="/etc/firstboot.applied"
ANYB="/etc/anyboot"

firstboot() {
    [ ! -d "$FIRSTB" ] && return
    [ ! -d "$FIRSTB_BACKUP" ] && mkdir -p $FIRSTB_BACKUP
    ls $FIRSTB/* > /dev/null 2>&1
    [ $? -ne 0 ] && return
    for f in $(ls $FIRSTB); do
        echo "[FirstBoot] Executing $FIRSTB/$f"
	cp -f $FIRSTB/$f $FIRSTB_BACKUP/
        rm -f $FIRSTB/$f # rm before executing or rebooting scripts will never get removed and loop forever
	$FIRSTB_BACKUP/$f
    done
}

anyboot() {
    [ ! -d "$ANYB" ] && return
    ls $ANYB/* > /dev/null 2>&1
    [ $? -ne 0 ] && return
    for f in $ANYB/*; do
        echo "[AnyBoot] Executing $f"
	$f
    done
}

start() {
    firstboot
    anyboot

    confine_node_enable
}


stop() {

    confine_daemon_stop

    confine_sliver_stop all

    if [ "$( uci -q get confine.node.state )" = "started" ] ; then
	confine.lib uci_set confine.node.state=applied
    fi
}
