# packages/confine/confine-system/files/etc/init.d/lxc
#!/bin/sh /etc/rc.common


START=80

#EXTRA_COMMANDS=""



start() {

    lxc.lib lxc_check

    local SL_NAME=

    for SL_NAME in $( lxc.lib uci_get_sections lxc container ); do

	if [ "$( uci -q get lxc.$SL_NAME.container_id || true )" = "${SL_NAME:0:2}" ] && [ "$( uci -q get lxc.$SL_NAME.auto_boot || true )" = "1" ]; then
	    lxc.lib lxc_start $SL_NAME
	fi
	
    done
}


stop() {

    lxc.lib lxc_check

    local SL_NAME=

    for SL_NAME in $( lxc-ls ); do

	local SL_TYPE="$( uci -q get lxc.$SL_NAME || true )"

	if [ "$SL_TYPE" = "container" ]; then
	    if lxc-info -n $SL_NAME 2>/dev/null | grep -e "^state:" | grep -e "RUNNING" >/dev/null ;then
	      	lxc.lib lxc_stop $SL_NAME
	    fi
	fi
    done

    lxc.lib loop_bridge_ifnames reset
}
