# files/etc/uci-defaults/confine-passwd.sh
#!/bin/sh

echo "Setting confine passwd (enabling ssh)"

passwd <<EOF
confine
confine
EOF


