# files/etc/rc.local
#!/bin/sh

mount -t cgroup none /cgroup

exit 0
