# packages/confine/dlep-app/files/dlep.init
#!/bin/sh /etc/rc.common

START=82
APP=/usr/sbin/dlep
ARGS="-l /etc/dlep.conf --set global.fork=true --set global.failfast=true"

start() {
	service_start $APP $ARGS
}

stop() {
	service_stop $APP $ARGS
}

restart() {
	stop
	sleep 1
	start
}
