# packages/confine/confine-system/files/etc/uci-defaults/confine-defer-uhttpd.sh
#!/bin/sh

echo "Disabling uhttpd (for deferred startup by confine.lib daemon)"

/etc/init.d/uhttpd disable

